# scripts/server/supervisor.conf
[program:tacnet]
directory = /home/tacnet-www/www/scripts/server/
user = tacnet-www
command = /home/tacnet-www/www/scripts/server/unicorn.sh
stdout_logfile = /home/tacnet-www/www/logs/supervisor.log
stderr_logfile = /home/tacnet-www/www/logs/supervisor_error.log

